# roles/install_node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
Documentation=https://prometheus.io/docs/guides/node-exporter/
After=network-online.target

[Service]
User={{ ansible_user }}
Restart=on-failure

ExecStart={{ install_path }}/node_exporter

[Install]
WantedBy=multi-user.target
